# files/root/service/logstash/run
#!/usr/bin/env sh
exec s6-setuidgid logstash /opt/logstash/bin/logstash -f /etc/logstash/logstash.conf 2>&1
