# files/root/service/elasticsearch/run
#!/usr/bin/env sh
exec s6-setuidgid elasticsearch /opt/elasticsearch/bin/elasticsearch 2>&1

