# files/root/service/kibana/run
#!/usr/bin/env sh

# reset permissions after docker volume remount
chown -R kibana:kibana /opt/kibana

exec s6-setuidgid kibana node --no-warnings /opt/kibana/src/cli 2>&1
