# files/root/service/traefik/run
#!/usr/bin/env sh

exec traefik --configFile=/opt/traefik/etc/traefik.toml 2>&1
